# chef/cookbooks/stunnel/templates/default/stunnel.socket.erb
# Dropped off by Chef

[Unit]
Description=Stunnel socket port=20200
After=network.target
Requires=network.target

[Socket]
ListenStream=20200
Accept=true

[Install]
WantedBy=containers.target
