# chef/cookbooks/stunnel/templates/default/stunnel.socket.erb
# Dropped off by Chef

[Unit]
Description=Stunnel socket
After=network.target

[Socket]
ListenStream=9999
