# chef/cookbooks/stunnel/templates/default/stunnel.service.erb
# Dropped off by Chef

[Unit]
Description=Stunnel service
After=network.target

[Service]
Type=forking
ExecStart=/usr/local/bin/stunnel /etc/stunnel/stunnel.conf
